# infra/backup/caffeine-level-critical-backup.service
[Unit]
Description=Caffeine Level Critical – nattlig backup (restic)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/caffeine-level-critical/infra/backup/backup.sh
